# root/etc/services.d/tftpd/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/in.tftpd --foreground --listen --secure $OPTS /data
